# docker/php/supervisor/conf.d/cron.ini
[program:cron]
stdout_logfile = /var/log/cron.log
stdout_logfile_maxbytes=0
stderr_logfile = /var/log/cron.error.log
stderr_logfile_maxbytes=0
autostart = false
autorestart = true
command=crond -f
stopsignal=TERM
