# docker/php/supervisor/conf.d/horizon.ini
[program:horizon]
directory=/app
command=/usr/local/bin/php artisan horizon
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor-horizon.log
autostart=false
autorestart=true
user=root
stopsignal=TERM
